# scripts/veil-xray-api.service
[Unit]
Description=Veil Xray API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/veil-v2ray
Environment="PATH=/root/veil-v2ray/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/root/veil-v2ray/venv/bin/python3 -m uvicorn api.main:app --host 127.0.0.1 --port 8000 --proxy-headers --forwarded-allow-ips 127.0.0.1 --timeout-keep-alive 10 --limit-concurrency 200 --backlog 2048
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=veil-xray-api

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ReadWritePaths=/usr/local/etc/xray
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true

[Install]
WantedBy=multi-user.target
